# init/carte.conf
# carte - carte jobfile

description "Commandline paste server"
author "Me <christophe@call-cc.be>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn
respawn limit 10 5

# Start the process
script
        /usr/local/carte -p 5003 -u "https://call-cc.be/pastes" -d /var/www/ssl/pastes
end script
